# Openfire/conf/etc/supervisor/conf.d/supervisor_openfire.conf
[supervisord]
nodaemon = true

[program:openfire]
command = /etc/init.d/openfire restart
user = root
auto_start = true
autorestart = true

stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

